# tradex.service
[Unit]
Description=TradeX Crypto Trading Bot - 24/7 Automation
After=network.target
Wants=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/TradeX
Environment=PATH=/home/ubuntu/TradeX/tradex_env/bin
Environment=PYTHONPATH=/home/ubuntu/TradeX
ExecStart=/home/ubuntu/TradeX/tradex_env/bin/python main.py --mode trade
Restart=always
RestartSec=10
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/TradeX

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradex

# Kill settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
